# provision/rec-stop.conf
[program:rec-stop]
command=/usr/local/bin/rec-stop.sh
priority=50
rutostart=true
autorestart=false
