# provision/rec-window.conf
[program:rec-window]
command=/usr/local/bin/rec-window.py
priority=30
autostart=true
autorestart=false
