# provision/rec-start.conf
[program:rec-start]
command=/usr/local/bin/rec-start.sh
priority=40
autostart=true
autorestart=true
