# provision/selenium-server.conf
[program:selenium-server]
command=java -jar /usr/share/selenium/selenium-server-standalone.jar
priority=20
autostart=true
autorestart=false
stopasgroup=true
