# provision/xvfb.conf
[program:xvfb]
command=/usr/bin/Xvfb -fp /usr/share/fonts/X11/misc/ :99 -screen 0 1280x1024x24
priority=10
autostart=true
autorestart=false
